# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:backend]
command=java -Xms256m -Xmx512m -jar /app/app.jar
directory=/app
autostart=true
autorestart=true
startsecs=30
startretries=3
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/backend_error.log
stderr_logfile_maxbytes=10MB
environment=JAVA_OPTS="-Xms256m -Xmx512m"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/nginx_error.log
stderr_logfile_maxbytes=10MB